# backup.timer
[Unit]
Description=Weekly backup

[Timer]
OnCalendar=*-*-1/3 18:00:00
Unit=backup.service

[Install]
WantedBy=multi-user.target
